# bin/hawck-inputd.service.in
[Unit]
Description=Hawck input daemon
#After=

[Service]
User=hawck-input
#Group=hawck-input
WorkingDirectory=/var/lib/hawck-input
#Type=forking
#PIDFile=/var/lib/hawck-input/pid
ExecStart=@PREFIX@/bin/hawck-inputd --no-fork
#ExecStop=/bin/bash -c 'kill $(cat /var/lib/hawck-input/pid)'

[Install]
WantedBy=multi-user.target

